# services/agent-hooks/ops/systemd/bloodbank-agent-hooks-health.service
[Unit]
Description=Bloodbank agent-hooks health → Holocene (validate deployed hook configs)
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/delorenj/code/33GOD/bloodbank
# `mise run health:hooks` runs the healthcheck under the repo toolchain (jsonschema
# present) and emits the snapshot to the Redis key Holocene reads. Sandbox only:
# our publishers run with BLOODBANK_ENABLED=false; foreign hooks are never executed.
ExecStart=/home/delorenj/.local/bin/mise run health:hooks
StandardOutput=append:/home/delorenj/code/33GOD/bloodbank/_bmad_output/evidence/health/timer.log
StandardError=append:/home/delorenj/code/33GOD/bloodbank/_bmad_output/evidence/health/timer.log
